# agents/com.syneira.scove-archive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.syneira.scove-archive</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sasha/V/VHome/agents/sync-cc-transcripts.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/sasha/V/VHome/logs/archive-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sasha/V/VHome/logs/archive-stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
